# supervisor/node.ini
[supervisord]
nodaemon=true
user=root

[program:node]
command=/usr/bin/start.sh
autostart=true
autorestart=true